# timekeeper-emo.service
[Unit]
Description=Timekeeper Emo-chan - NFC Time Tracking with BOCCO emo
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/dev/timekeeper-emo-chan
Environment="PATH=/home/pi/dev/timekeeper-emo-chan/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
ExecStart=/home/pi/dev/timekeeper-emo-chan/venv/bin/python main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Graceful shutdown
TimeoutStopSec=10
KillMode=mixed
KillSignal=SIGINT

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
